# newvk-bot.service
[Unit]
Description=NewVK VK webhook bot
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/newvk
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/mkdir -p /root/newvk/logs
ExecStart=/usr/local/bin/python3 /root/newvk/vk_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
